# kindle/scripts/dashboard.conf
#kate: syntax bash;
description "Kindle Dashboard"
#version "v0.1"

start on started pillow and started cmd and started lab126 and started framework

script
    su root -c "sleep 150 && sh /mnt/us/scripts/dashboardscript.sh"
    return 0
end script

post-stop script
    return 0
end script